# scripts/retroi-failure@.service
[Unit]
Description=Failure handler for %i

[Service]
Type=oneshot
ExecStart=/bin/bash -c '$RETROI_DIR/scripts/update_project.sh main'
